# kittentts/rootfs/etc/s6-overlay/s6-rc.d/kittentts/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "KittenTTS stopped"
